# etc/supervisor/conf.d/10-xvfb.conf
[program:xvfb]
command = /usr/bin/Xvfb :1 -screen 0 %(ENV_XRES)s
priority=10
