# etc/supervisor/conf.d/40-novnc.conf
[program:novnc]
priority=40
directory=/usr/share/novnc/
command=/bin/bash /usr/share/novnc/utils/novnc_proxy --web /usr/share/novnc --listen 6080
stopasgroup=true
